# common/local.newapi.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>local.newapi.service</string> <!-- 服务的唯一标签 -->

    <key>ProgramArguments</key>
    <array>
        <string>/Users/liujiahao/WorkBench/new-api/newapi</string> <!-- 你的 Go 服务可执行文件的绝对路径 -->
        <!-- 如果你的 Go 服务需要命令行参数，在这里添加，例如： -->
        <!-- <string>--config</string> -->
        <!-- <string>/Users/liujiahao/WorkBench/new-api/config.json</string> -->
    </array>

    <key>RunAtLoad</key>
    <true/> <!-- 表示在加载此 .plist 文件时立即启动服务 -->

    <key>KeepAlive</key>
    <true/> <!-- 表示如果服务意外终止，launchd 会尝试重新启动它 -->

    <key>WorkingDirectory</key>
    <string>/Users/liujiahao/WorkBench/new-api/</string> <!-- 设置服务的工作目录 -->

    <key>StandardOutPath</key>
    <string>/Users/liujiahao/WorkBench/new-api/newapi.service.out.log</string> <!-- (可选) 标准输出日志文件路径 -->

    <key>StandardErrorPath</key>
    <string>/Users/liujiahao/WorkBench/new-api/newapi.service.out.log</string> <!-- (可选) 标准错误日志文件路径 -->

    <!-- 如果你希望服务在当前用户登录时启动 (放在 ~/Library/LaunchAgents/) -->
    <!-- 并且以当前用户身份运行，通常不需要显式设置 UserName -->

    <!-- 如果你希望
